# lib/udev/rules.d/999-android-system.rules
# Generated by the Mer uevent to udev rule generator
# in the droid-hal-* package
#
# This file contains udev rules derived from the uevent.rc
# The original rc file contents are included preceded by #:
#
#
################
# From ueventd.rc

#:subsystem adf
# WARNING rule in comment above was not converted
#:    devname uevent_devname
# WARNING rule in comment above was not converted
#:
#:subsystem graphics
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/graphics
# WARNING rule in comment above was not converted
#:
#:subsystem drm
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/dri
# WARNING rule in comment above was not converted
#:
#:subsystem oncrpc
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/oncrpc
# WARNING rule in comment above was not converted
#:
#:subsystem adsp
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/adsp
# WARNING rule in comment above was not converted
#:
#:subsystem msm_camera
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/msm_camera
# WARNING rule in comment above was not converted
#:
#:subsystem input
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/input
# WARNING rule in comment above was not converted
#:
#:subsystem mtd
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/mtd
# WARNING rule in comment above was not converted
#:
#:subsystem sound
# WARNING rule in comment above was not converted
#:    devname uevent_devpath
# WARNING rule in comment above was not converted
#:    dirname /dev/snd
# WARNING rule in comment above was not converted
#:
#:# ueventd can only set permissions on device nodes and their associated
#:# sysfs attributes, not on arbitrary paths.
#:#
#:# format for /dev rules: devname mode uid gid
#:# format for /sys rules: nodename attr mode uid gid
#:# shortcut: "mtd@NN" expands to "/dev/mtd/mtdNN"
#:
#:/dev/null                 0666   root       root
KERNEL=="null", MODE="0666", GROUP="root", OWNER="root"
#:/dev/zero                 0666   root       root
KERNEL=="zero", MODE="0666", GROUP="root", OWNER="root"
#:/dev/full                 0666   root       root
KERNEL=="full", MODE="0666", GROUP="root", OWNER="root"
#:/dev/ptmx                 0666   root       root
KERNEL=="ptmx", MODE="0666", GROUP="root", OWNER="root"
#:/dev/tty                  0666   root       root
KERNEL=="tty", MODE="0666", GROUP="root", OWNER="root"
#:/dev/random               0666   root       root
KERNEL=="random", MODE="0666", GROUP="root", OWNER="root"
#:/dev/urandom              0666   root       root
KERNEL=="urandom", MODE="0666", GROUP="root", OWNER="root"
#:# Make HW RNG readable by group system to let EntropyMixer read it.
#:/dev/hw_random            0440   root       system
KERNEL=="hw_random", MODE="0440", GROUP="system", OWNER="root"
#:/dev/ashmem               0666   root       root
KERNEL=="ashmem", MODE="0666", GROUP="root", OWNER="root"
#:/dev/binder               0666   root       root
KERNEL=="binder", MODE="0666", GROUP="root", OWNER="root"
#:/dev/hwbinder             0666   root       root
KERNEL=="hwbinder", MODE="0666", GROUP="root", OWNER="root"
#:/dev/vndbinder            0666   root       root
KERNEL=="vndbinder", MODE="0666", GROUP="root", OWNER="root"
#:
#:/dev/pmsg0                0222   root       log
KERNEL=="pmsg0", MODE="0222", GROUP="log", OWNER="root"
#:
#:# the msm hw3d client device node is world writable/readable.
#:/dev/msm_hw3dc            0666   root       root
KERNEL=="msm_hw3dc", MODE="0666", GROUP="root", OWNER="root"
#:
#:# gpu driver for adreno200 is globally accessible
#:/dev/kgsl                 0666   root       root
KERNEL=="kgsl", MODE="0666", GROUP="root", OWNER="root"
#:
#:# kms driver for drm based gpu
#:/dev/dri/*                0666   root       graphics
# Known droid dir drm/ is for subsystem dri
SYMLINK=="dri/*", MODE="0666", GROUP="graphics", OWNER="root"
SUBSYSTEM=="drm", KERNEL=="*", MODE="0666", GROUP="graphics", OWNER="root"
#:
#:# these should not be world writable
#:/dev/diag                 0660   radio      radio
KERNEL=="diag", MODE="0660", GROUP="radio", OWNER="radio"
#:/dev/diag_arm9            0660   radio      radio
KERNEL=="diag_arm9", MODE="0660", GROUP="radio", OWNER="radio"
#:/dev/ttyMSM0              0600   bluetooth  bluetooth
KERNEL=="ttyMSM0", MODE="0600", GROUP="bluetooth", OWNER="bluetooth"
#:/dev/uhid                 0660   uhid       uhid
KERNEL=="uhid", MODE="0660", GROUP="uhid", OWNER="uhid"
#:/dev/uinput               0660   system     bluetooth
KERNEL=="uinput", MODE="0660", GROUP="bluetooth", OWNER="system"
#:/dev/alarm                0664   system     radio
KERNEL=="alarm", MODE="0664", GROUP="radio", OWNER="system"
#:/dev/rtc0                 0640   system     system
KERNEL=="rtc0", MODE="0640", GROUP="system", OWNER="system"
#:/dev/tty0                 0660   root       system
KERNEL=="tty0", MODE="0660", GROUP="system", OWNER="root"
#:/dev/graphics/*           0660   root       graphics
# Known droid dir graphics/ is for subsystem graphics
SYMLINK=="graphics/*", MODE="0660", GROUP="graphics", OWNER="root"
SUBSYSTEM=="graphics", KERNEL=="*", MODE="0660", GROUP="graphics", OWNER="root"
#:/dev/msm_hw3dm            0660   system     graphics
KERNEL=="msm_hw3dm", MODE="0660", GROUP="graphics", OWNER="system"
#:/dev/input/*              0660   root       input
#:/dev/v4l-touch*           0660   root       input
KERNEL=="v4l-touch*", MODE="0660", GROUP="input", OWNER="root"
#:/dev/eac                  0660   root       audio
KERNEL=="eac", MODE="0660", GROUP="audio", OWNER="root"
#:/dev/cam                  0660   root       camera
KERNEL=="cam", MODE="0660", GROUP="camera", OWNER="root"
#:/dev/pmem                 0660   system     graphics
KERNEL=="pmem", MODE="0660", GROUP="graphics", OWNER="system"
#:/dev/pmem_adsp*           0660   system     audio
KERNEL=="pmem_adsp*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/pmem_camera*         0660   system     camera
KERNEL=="pmem_camera*", MODE="0660", GROUP="camera", OWNER="system"
#:/dev/oncrpc/*             0660   root       system
# Known droid dir oncrpc/ is for subsystem oncrpc
SYMLINK=="oncrpc/*", MODE="0660", GROUP="system", OWNER="root"
SUBSYSTEM=="oncrpc", KERNEL=="*", MODE="0660", GROUP="system", OWNER="root"
#:/dev/adsp/*               0660   system     audio
# Known droid dir adsp/ is for subsystem adsp
SYMLINK=="adsp/*", MODE="0660", GROUP="audio", OWNER="system"
SUBSYSTEM=="adsp", KERNEL=="*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/snd/*                0660   system     audio
# Standard subsystem rule
SUBSYSTEM=="snd", KERNEL=="*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/mt9t013              0660   system     system
KERNEL=="mt9t013", MODE="0660", GROUP="system", OWNER="system"
#:/dev/msm_camera/*         0660   system     system
# Known droid dir msm_camera/ is for subsystem msm_camera
SYMLINK=="msm_camera/*", MODE="0660", GROUP="system", OWNER="system"
SUBSYSTEM=="msm_camera", KERNEL=="*", MODE="0660", GROUP="system", OWNER="system"
#:/dev/akm8976_daemon       0640   compass    system
KERNEL=="akm8976_daemon", MODE="0640", GROUP="system", OWNER="compass"
#:/dev/akm8976_aot          0640   compass    system
KERNEL=="akm8976_aot", MODE="0640", GROUP="system", OWNER="compass"
#:/dev/akm8973_daemon       0640   compass    system
KERNEL=="akm8973_daemon", MODE="0640", GROUP="system", OWNER="compass"
#:/dev/akm8973_aot          0640   compass    system
KERNEL=="akm8973_aot", MODE="0640", GROUP="system", OWNER="compass"
#:/dev/bma150               0640   compass    system
KERNEL=="bma150", MODE="0640", GROUP="system", OWNER="compass"
#:/dev/cm3602               0640   compass    system
KERNEL=="cm3602", MODE="0640", GROUP="system", OWNER="compass"
#:/dev/akm8976_pffd         0640   compass    system
KERNEL=="akm8976_pffd", MODE="0640", GROUP="system", OWNER="compass"
#:/dev/lightsensor          0640   system     system
KERNEL=="lightsensor", MODE="0640", GROUP="system", OWNER="system"
#:/dev/msm_pcm_out*         0660   system     audio
KERNEL=="msm_pcm_out*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/msm_pcm_in*          0660   system     audio
KERNEL=="msm_pcm_in*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/msm_pcm_ctl*         0660   system     audio
KERNEL=="msm_pcm_ctl*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/msm_snd*             0660   system     audio
KERNEL=="msm_snd*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/msm_mp3*             0660   system     audio
KERNEL=="msm_mp3*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/audience_a1026*      0660   system     audio
KERNEL=="audience_a1026*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/tpa2018d1*           0660   system     audio
KERNEL=="tpa2018d1*", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/msm_audpre           0660   system     audio
KERNEL=="msm_audpre", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/msm_audio_ctl        0660   system     audio
KERNEL=="msm_audio_ctl", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/htc-acoustic         0660   system     audio
KERNEL=="htc-acoustic", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/vdec                 0660   system     audio
KERNEL=="vdec", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/q6venc               0660   system     audio
KERNEL=="q6venc", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/snd/dsp              0660   system     audio
# Standard subsystem rule
SUBSYSTEM=="snd", KERNEL=="dsp", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/snd/dsp1             0660   system     audio
# Standard subsystem rule
SUBSYSTEM=="snd", KERNEL=="dsp1", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/snd/mixer            0660   system     audio
# Standard subsystem rule
SUBSYSTEM=="snd", KERNEL=="mixer", MODE="0660", GROUP="audio", OWNER="system"
#:/dev/smd0                 0640   radio      radio
KERNEL=="smd0", MODE="0640", GROUP="radio", OWNER="radio"
#:/dev/qmi                  0640   radio      radio
KERNEL=="qmi", MODE="0640", GROUP="radio", OWNER="radio"
#:/dev/qmi0                 0640   radio      radio
KERNEL=="qmi0", MODE="0640", GROUP="radio", OWNER="radio"
#:/dev/qmi1                 0640   radio      radio
KERNEL=="qmi1", MODE="0640", GROUP="radio", OWNER="radio"
#:/dev/qmi2                 0640   radio      radio
KERNEL=="qmi2", MODE="0640", GROUP="radio", OWNER="radio"
#:/dev/bus/usb/*            0660   root       usb
# usb rule
SUBSYSTEM=="usb", KERNEL=="*", MODE="0660", GROUP="usb", OWNER="root"
#:/dev/mtp_usb              0660   root       mtp
KERNEL=="mtp_usb", MODE="0660", GROUP="mtp", OWNER="root"
#:/dev/usb_accessory        0660   root       usb
KERNEL=="usb_accessory", MODE="0660", GROUP="usb", OWNER="root"
#:/dev/tun                  0660   system     vpn
KERNEL=="tun", MODE="0660", GROUP="vpn", OWNER="system"
#:
#:# CDMA radio interface MUX
#:/dev/ts0710mux*           0640   radio      radio
KERNEL=="ts0710mux*", MODE="0640", GROUP="radio", OWNER="radio"
#:/dev/ppp                  0660   radio      vpn
KERNEL=="ppp", MODE="0660", GROUP="vpn", OWNER="radio"
#:
#:# sysfs properties
#:/sys/devices/platform/trusty.*      trusty_version        0440  root   log
# sys rule
DEVPATH=="/devices/platform/trusty.*", RUN+="/bin/chmod 0440 /sys/$devpath/trusty_version", RUN+="/bin/chown root:log /sys/$devpath/trusty_version"
#:/sys/devices/virtual/input/input*   enable      0660  root   input
# sys rule
DEVPATH=="/devices/virtual/input/input*", RUN+="/bin/chmod 0660 /sys/$devpath/enable", RUN+="/bin/chown root:input /sys/$devpath/enable"
#:/sys/devices/virtual/input/input*   poll_delay  0660  root   input
# sys rule
DEVPATH=="/devices/virtual/input/input*", RUN+="/bin/chmod 0660 /sys/$devpath/poll_delay", RUN+="/bin/chown root:input /sys/$devpath/poll_delay"
#:/sys/devices/virtual/usb_composite/*   enable      0664  root   system
# sys rule
DEVPATH=="/devices/virtual/usb_composite/*", RUN+="/bin/chmod 0664 /sys/$devpath/enable", RUN+="/bin/chown root:system /sys/$devpath/enable"
#:/sys/devices/system/cpu/cpu*   cpufreq/scaling_max_freq   0664  system system
# sys rule
DEVPATH=="/devices/system/cpu/cpu*", RUN+="/bin/chmod 0664 /sys/$devpath/cpufreq/scaling_max_freq", RUN+="/bin/chown system:system /sys/$devpath/cpufreq/scaling_max_freq"
#:/sys/devices/system/cpu/cpu*   cpufreq/scaling_min_freq   0664  system system
# sys rule
DEVPATH=="/devices/system/cpu/cpu*", RUN+="/bin/chmod 0664 /sys/$devpath/cpufreq/scaling_min_freq", RUN+="/bin/chown system:system /sys/$devpath/cpufreq/scaling_min_freq"
#:
#:# DVB API device nodes
#:/dev/dvb*                 0660   root       system
KERNEL=="dvb*", MODE="0660", GROUP="system", OWNER="root"
################
